# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Stop and disable the old service
        if [ -d /run/systemd/system ]; then
            systemctl stop halos-mdns-publisher.service 2>/dev/null || true
            systemctl disable halos-mdns-publisher.service 2>/dev/null || true
            systemctl daemon-reload
        fi

        # Remove the old binary and service file
        rm -f /usr/bin/halos-mdns-publisher
        rm -f /etc/systemd/system/halos-mdns-publisher.service
        rm -f /lib/systemd/system/halos-mdns-publisher.service
        ;;
esac

exit 0
